# test-infra/ansible/roles/node_identifier/files/node-identifier.service
[Unit]
Description=One-shot Node Identifier and Runner Registration
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/identify_and_register.sh
EnvironmentFile=-/etc/runner-config.env
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
